# watchdogd.service.in
[Unit]
Description=Advanced watchdog daemon
Documentation=man:watchdogd(8) man:watchdogctl(1) man:watchdogd.conf(5)

[Service]
EnvironmentFile=-@SYSCONFDIR@/default/watchdogd
ExecStart=@SBINDIR@/watchdogd -ns $WATCHDOGD_OPTS
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
