# watchdogd.service.in
[Unit]
Description=Advanced watchdog daemon
Documentation=man:watchdogd(8) man:watchdogctl(1)
Documentation=file:@DOCDIR@/README.md
Documentation=http://github.com/troglobit/watchdogd
Documentation=http://troglobit.com/project/watchdogd

[Service]
Type=simple
Restart=always
RestartSec=3
ExecStart=@SBINDIR@/watchdogd -n -s -a 0.9 -f 0.9 -m 0.9

[Install]
WantedBy=multi-user.target
